# etc/systemd/system/sysgit-autocommit.service
[Unit]
Description=Sysgit autocommit snapshot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sysgit -autocommit
