# etc/systemd/system/sysgit-autocommit.timer
[Unit]
Description=Nightly sysgit autocommit snapshot

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
